# server/services/zfs-scrub/zfs-scrub.service
[Unit]
Description=ZFS pool scrub

[Service]
Type=oneshot
ExecStart=/sbin/zpool scrub zmain

