# server/services/zfs-report/zfs-report.service
[Unit]
Description=ZFS Subsystem Report
After=zfs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/genzfsreport