# server/services/full-backup/full-backup.timer
[Unit]
Description=Weekly system backup

[Timer]
Persistent=true
OnCalendar=Sat, 05:30
AccuracySec=10min

[Install]
WantedBy=multi-user.target
