# server/services/full-backup/full-backup.service
[Unit]
Description=Full system backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup /mnt/backup/archives/server

